# scripts/Scripts/com.user.mac.maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.mac.maintenance</string>

  <!-- Run the script through zsh so $HOME expands inside the command -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>$HOME/Scripts/mac_maintenance/maintenance.sh --auto</string>
  </array>

  <!-- Example: every day at 03:00. Adjust as you like or remove to trigger manually. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>$HOME/Library/Logs/mac_maintenance.launchd.out</string>
  <key>StandardErrorPath</key>
  <string>$HOME/Library/Logs/mac_maintenance.launchd.err</string>

  <!-- Don't keep the job alive if it exits; it will run based on the schedule. -->
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
